# stress-sys.service
[Unit]
Description=Stress test the system forever

[Service]
Type=idle
ExecStart=/opt/stress-app/stress-sys.sh

[Install]
WantedBy=multi-user.target
